# docker/root/etc/services.d/pical/run
#!/usr/bin/with-contenv sh

/app/pical.sh
